# src/main/dist/moasdawiki.service
[Unit]
Description=MoasdaWiki

[Service]
Type=forking
TimeoutStartSec=0
ExecStart=/path/to/moasdawiki.sh start
ExecStop=/path/to/moasdawiki.sh stop
PIDFile=/var/run/moasdawiki.pid

[Install]
WantedBy=multi-user.target
